# deploy/ansible/roles/icaro/files/etc/systemd/system/caddy.service.d/override.conf
[Unit]
After=cloud-final.service

[Service]
EnvironmentFile=-/etc/caddy/envfile
Restart=on-failure
RestartSec=5
